# weekly-vm-backup.cron
# Put this file to /etc/cron.d
# Start weekly VM backup

SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin

# CHANGE THIS!
MAILFROM=robot@firm.com
MAILTO=admin@firm.com

@weekly root /root/weekly-vm-backup
